# deploy/leetcode-senpai-testing-sync.timer
[Unit]
Description=Periodically pin leetcode-senpai testing to the integration head

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Unit=leetcode-senpai-testing-sync.service

[Install]
WantedBy=timers.target
